# nexus-ai-addon/rootfs/etc/services.d/nexus-ai/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Nexus AI service
# ==============================================================================
bashio::log.info "Starting Nexus AI..."

cd /app
exec /app/run.sh
